# system/claw-dashboard.service
[Unit]
Description=Progress Claw OS machine dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=araya
Group=araya
WorkingDirectory=/home/araya/Projects/Progress-Claw-OS
Environment=PYTHONUNBUFFERED=1
Environment=CLAW_ARDUINO_TRANSPORT=gpio
ExecStart=/home/araya/Projects/Progress-Claw-OS/.venv/bin/python /home/araya/Projects/Progress-Claw-OS/main.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
